# brain/infrastructure/config/supervisord.d/brain_google_api.conf
# Brain Google API Service Configuration for supervisord
# MCP server for Gmail / Drive / Calendar / Docs / Sheets / Tasks / People

[program:brain_google_api]
command=/brain/infrastructure/service/brain_google_api/releases/v0.1.0/bin/service-brain_google_api
directory=/brain/infrastructure/service/brain_google_api
environment=GOOGLE_CREDENTIALS_PATH="/brain/secrets/brain_google_api/credentials.json",DAEMON_SOCKET="/tmp/brain_ipc.sock",SERVICE_NAME="service-brain_google_api"
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/brain/runtime/logs/brain_google_api.log
stderr_logfile=/brain/runtime/logs/brain_google_api_error.log
