# brain/infrastructure/config/supervisord.d/brain_dashboard.conf
# Brain Dashboard Service Configuration for supervisord
# BS-025 brain_dashboard v2.2.0

[program:brain_dashboard]
command=/opt/venv/bin/python -m uvicorn app:app --host 0.0.0.0 --port 8080 --log-level info
directory=/brain/infrastructure/service/brain_dashboard/releases/current/src
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/xkagent_infra/runtime/logs/brain_dashboard.log
stderr_logfile=/xkagent_infra/runtime/logs/brain_dashboard_error.log
environment=DASHBOARD_VERSION="2.2.0",DASHBOARD_SERVICE="service-brain_dashboard",PYTHONPATH="/brain/infrastructure/service/utils/ipc/src"
