# brain/infrastructure/config/supervisord.d/brain_task_manager.conf
# Brain Task Manager Service Configuration for supervisord
# BS-026 brain_task_manager

[program:brain_task_manager]
command=/brain/infrastructure/service/brain_task_manager/bin/brain_task_manager --config /brain/infrastructure/config/brain_task_manager/brain_task_manager.json
directory=/brain/infrastructure/service/brain_task_manager
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/brain/runtime/logs/brain_task_manager.log
stderr_logfile=/brain/runtime/logs/brain_task_manager_error.log
