# brain/infrastructure/config/supervisord.d/health_check.conf
# Brain Health Check Service Configuration for supervisord
# BS-024-health-check v1.0.0

[program:health_check]
command=python3 /brain/infrastructure/service/brain_health_check/health_server.py
directory=/brain/infrastructure/service/brain_health_check
autostart=true
autorestart=true
startretries=3
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/xkagent_infra/runtime/logs/health_check.out.log
stderr_logfile=/xkagent_infra/runtime/logs/health_check.err.log
environment=BRAIN_HEALTH_PORT="8765",BRAIN_IPC_SOCKET_PATH="/tmp/brain_ipc.sock",HEALTH_CHECK_TIMEOUT="2.0",HEALTH_CHECK_PROBE_TIMEOUT="0.5",HEALTH_CHECK_VERSION="1.0.0",DAEMON_SOCKET="/tmp/brain_ipc.sock",SERVICE_NAME="service-brain_health_check"
