# brain/infrastructure/config/supervisord.d/telegram_services.conf
# Telegram Services Configuration for supervisord
# Single process mode

# Telegram Service (registers as service-telegram_api)
[program:service_telegram_api]
command=/brain/infrastructure/service/brain_telegram_api/releases/current/bin/telegram_api_service
directory=/brain/infrastructure/service/brain_telegram_api
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
stdout_logfile=/xkagent_infra/runtime/logs/telegram_api.log
stderr_logfile=/xkagent_infra/runtime/logs/telegram_api_error.log
environment=
  TELEGRAM_BOT_TOKEN_FILE=/brain/secrets/IM/telegram/bot_token_2.env,
  TELEGRAM_CONFIG=/brain/infrastructure/config/third_api/telegram/telegram.yaml,
  DAEMON_SOCKET=/tmp/brain_ipc.sock,
  SERVICE_NAME=service-telegram_api
