# groups/brain/platform/docker/configs/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
user=root
logfile=/xkagent_infra/runtime/logs/supervisord.log
pidfile=/var/run/supervisord.pid
; Common environment inherited by all child programs
environment=HOME="/root",BRAIN_PATH="/xkagent_infra/brain",BRAIN_IPC_SOCKET="/tmp/brain_ipc.sock",DAEMON_SOCKET="/tmp/brain_ipc.sock",PATH="/root/.local/bin:/opt/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/xkagent_infra/runtime/logs/sshd.out.log
stderr_logfile=/xkagent_infra/runtime/logs/sshd.err.log

[program:brain_ipc]
command=/xkagent_infra/brain/infrastructure/service/brain_ipc/bin/brain_ipc
directory=/xkagent_infra/brain/infrastructure/service/brain_ipc
autostart=true
autorestart=true
priority=10
stopasgroup=true
killasgroup=true
environment=BRAIN_TMUX_NOTIFY="1"
stdout_logfile=/xkagent_infra/runtime/logs/brain_ipc.out.log
stderr_logfile=/xkagent_infra/runtime/logs/brain_ipc.err.log

[program:agent_orchestrator]
command=/bin/bash -lc 'targets="${AGENTCTL_AUTOSTART_AGENTS:-agent-brain_manager}"; targets="${targets//,/ }"; compact="${targets// /}"; if [ -n "${compact}" ]; then python3 bin/brain-agentctl start --apply ${targets} || true; fi; exec python3 bin/brain-agentctl serve'
directory=/xkagent_infra/brain/infrastructure/service/agentctl
autostart=true
autorestart=true
priority=20
environment=AGENTCTL_CONFIG_DIR="/xkagent_infra/brain/infrastructure/config/agentctl",AGENT_MANAGER_CONFIG_DIR="/xkagent_infra/brain/infrastructure/config/agentctl",PYTHONPATH="/xkagent_infra/brain/infrastructure/service/agentctl"
stdout_logfile=/xkagent_infra/runtime/logs/agent_orchestrator.out.log
stderr_logfile=/xkagent_infra/runtime/logs/agent_orchestrator.err.log

[program:brain_timer]
command=/xkagent_infra/brain/infrastructure/service/brain_timer/bin/service_timer
directory=/xkagent_infra/brain/infrastructure/service/brain_timer
autostart=true
autorestart=true
priority=22
startsecs=3
startretries=5
stopasgroup=true
killasgroup=true
environment=PYTHONPATH="/xkagent_infra/groups/brain/projects/infrastructure/brain_timer/src"
stdout_logfile=/xkagent_infra/runtime/logs/brain_timer.supervisor.log
stderr_logfile=/xkagent_infra/runtime/logs/brain_timer.supervisor.err.log

[program:brain_task_manager]
command=/xkagent_infra/brain/infrastructure/service/brain_task_manager/bin/brain_task_manager --config /xkagent_infra/brain/infrastructure/config/brain_task_manager/brain_task_manager.json
directory=/xkagent_infra/brain/infrastructure/service/brain_task_manager
autostart=true
autorestart=true
priority=23
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/xkagent_infra/runtime/logs/brain_task_manager.supervisor.log
stderr_logfile=/xkagent_infra/runtime/logs/brain_task_manager.supervisor.err.log

[program:brain_health_check]
command=/xkagent_infra/brain/infrastructure/service/brain_health_check/bin/health_server.py
directory=/xkagent_infra/brain/infrastructure/service/brain_health_check
autostart=true
autorestart=true
priority=24
startsecs=3
startretries=3
stopasgroup=true
killasgroup=true
environment=SERVICE_NAME="service-brain_health_check",BRAIN_HEALTH_PORT="8765"
stdout_logfile=/xkagent_infra/runtime/logs/brain_health_check.supervisor.log
stderr_logfile=/xkagent_infra/runtime/logs/brain_health_check.supervisor.err.log

[program:service_telegram_api]
command=/xkagent_infra/brain/infrastructure/service/brain_telegram_api/bin/telegram_api_service
directory=/xkagent_infra/brain/infrastructure/service/brain_telegram_api
autostart=true
autorestart=true
priority=25
startsecs=3
startretries=3
stopasgroup=true
killasgroup=true
environment=PYTHONPATH="/xkagent_infra/groups/brain/projects/infrastructure/brain_telegram_api/src"
stdout_logfile=/xkagent_infra/runtime/logs/telegram_api.supervisor.log
stderr_logfile=/xkagent_infra/runtime/logs/telegram_api.supervisor.err.log

[program:brain_google_api]
command=/xkagent_infra/brain/infrastructure/service/brain_google_api/bin/service-brain_google_api
directory=/xkagent_infra/brain/infrastructure/service/brain_google_api
autostart=true
autorestart=true
priority=26
startsecs=3
startretries=3
stopasgroup=true
killasgroup=true
environment=SERVICE_NAME="service-brain_google_api",GOOGLE_CREDENTIALS_PATH="/brain/secrets/brain_google_api/credentials.json"
stdout_logfile=/xkagent_infra/runtime/logs/brain_google_api.supervisor.log
stderr_logfile=/xkagent_infra/runtime/logs/brain_google_api.supervisor.err.log

[program:brain_agent_proxy]
command=/bin/bash -lc 'exec /xkagent_infra/brain/infrastructure/service/brain_agent_proxy/bin/brain_agent_proxy run --port 8210 --host 0.0.0.0 --log-level info'
directory=/xkagent_infra/brain/infrastructure/service/brain_agent_proxy
autostart=true
autorestart=true
priority=27
startsecs=3
startretries=5
stopasgroup=true
killasgroup=true
environment=SERVICE_NAME="service-brain_agent_proxy"
stdout_logfile=/xkagent_infra/runtime/logs/brain_agent_proxy.supervisor.log
stderr_logfile=/xkagent_infra/runtime/logs/brain_agent_proxy.supervisor.err.log

[program:brain_gateway]
command=/xkagent_infra/brain/infrastructure/service/brain_gateway/bin/brain_gateway --config /xkagent_infra/brain/infrastructure/service/brain_gateway/config/brain_gateway.json
directory=/xkagent_infra/brain/infrastructure/service/brain_gateway
autostart=true
autorestart=true
priority=30
startsecs=3
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/xkagent_infra/runtime/logs/brain_gateway.supervisor.log
stderr_logfile=/xkagent_infra/runtime/logs/brain_gateway.supervisor.err.log

[program:brain_dashboard]
command=/opt/venv/bin/python -m uvicorn app:app --host 0.0.0.0 --port 8080 --log-level info
directory=/xkagent_infra/brain/infrastructure/service/brain_dashboard/current/src
autostart=true
autorestart=true
priority=28
startsecs=3
startretries=3
stopasgroup=true
killasgroup=true
environment=DASHBOARD_VERSION="2.2.0",DASHBOARD_SERVICE="service-brain_dashboard"
stdout_logfile=/xkagent_infra/runtime/logs/brain_dashboard.supervisor.log
stderr_logfile=/xkagent_infra/runtime/logs/brain_dashboard.supervisor.err.log
