# brain/platform/docker/configs/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/var/log/sshd.err.log
stdout_logfile=/var/log/sshd.out.log

[program:brain_ipc]
command=/xkagent_infra/brain/infrastructure/service/brain_ipc/bin/current/brain_ipc
directory=/xkagent_infra/brain/infrastructure/service/brain_ipc
autostart=true
autorestart=true
priority=10
environment=HOME="/root",BRAIN_PATH="/xkagent_infra/brain",BRAIN_IPC_SOCKET="/tmp/brain_ipc.sock",BRAIN_TMUX_NOTIFY="1"
stderr_logfile=/var/log/brain_ipc.err.log
stdout_logfile=/var/log/brain_ipc.out.log

[program:agent_orchestrator]
command=/bin/bash -lc 'python3 bin/brain-agentctl start --apply all || true; exec python3 bin/brain-agentctl serve'
directory=/xkagent_infra/brain/infrastructure/service/agentctl
autostart=true
autorestart=true
priority=20
environment=HOME="/root",BRAIN_PATH="/xkagent_infra/brain",BRAIN_IPC_SOCKET="/tmp/brain_ipc.sock",AGENTCTL_CONFIG_DIR="/xkagent_infra/brain/infrastructure/config/agentctl",AGENT_MANAGER_CONFIG_DIR="/xkagent_infra/brain/infrastructure/config/agentctl",PYTHONPATH="/xkagent_infra/brain/infrastructure/service/agentctl",PATH="/root/.local/bin:/opt/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
stderr_logfile=/var/log/agent_orchestrator.err.log
stdout_logfile=/var/log/agent_orchestrator.out.log
